# systemd/cbos.service
[Unit]
Description=CBOS - Claude Code Session Manager API
After=network.target

[Service]
Type=simple
User=bisenbek
Group=bisenbek
WorkingDirectory=/home/bisenbek/projects/nominate/cbos
# PATH includes ~/.local/bin for claude CLI
Environment="PATH=/home/bisenbek/.local/bin:/home/bisenbek/.pyenv/versions/nominates/bin:/usr/local/bin:/usr/bin:/bin"
Environment="CBOS_LOG_LEVEL=INFO"
ExecStart=/home/bisenbek/.pyenv/versions/nominates/bin/uvicorn cbos.api.main:app --host 127.0.0.1 --port 32205
Restart=always
RestartSec=5
SyslogIdentifier=cbos
StandardOutput=journal
StandardError=journal

# Minimal security hardening
# Note: Removed ProtectHome=read-only as it interferes with streaming file watching
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
